# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:sunshine]
command=/usr/bin/sunshine
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sunshine.log
stderr_logfile=/var/log/supervisor/sunshine_error.log
environment=DISPLAY=:0,SUNSHINE_USERNAME=admin,SUNSHINE_PASSWORD=change_me

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1920x1080x24
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb_error.log
priority=1

[program:ros2-launch]
command=/bin/bash -c "source /opt/ros/jazzy/setup.bash && source /home/ros/ros2_ws/install/setup.bash 2>/dev/null; exec ros2 launch ros2_launch.py"
user=root
directory=/home/ros/ros2_ws
autostart=false
autorestart=true
stdout_logfile=/var/log/supervisor/ros2_launch.log
stderr_logfile=/var/log/supervisor/ros2_launch_error.log
environment=ROS_DOMAIN_ID=0,DISPLAY=:0

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
